# .config/systemd/user/rclone-sync.service
[Unit]
Description=Sync files between different remotes via rclone
[Service]
Type=simple
ExecStart=/home/jjk/bin/rclone-sync
