# containers/srb2/home/srb2/.config/systemd/systemd/user/srb2_chars.timer
[Unit]
Description=Start srb2 chars server instead of the normal one

[Timer]
OnCalendar=Fri 00:00
Persistent=true

[Install]
WantedBy=timers.target
